# system76-extension/cursor-agent.service
[Unit]
Description=Super Sovereign Cursor-Agent CLI - System76 Extension
Documentation=https://github.com/tig08bitties/Agent
After=network.target system76-power.service
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/cbrew system76-activate
ExecStop=/usr/local/bin/cbrew system76-deactivate
StandardOutput=journal
StandardError=journal
User=root
Environment="CURSOR_AGENT_HOME=/root/cursor-agent"
Environment="NODE_PATH=/root/cursor-agent/node_modules"

[Install]
WantedBy=multi-user.target
